# personal-garmin-ai-coach.service
[Unit]
Description=Personal Garmin AI Coach
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/PersonalGarminAICoach
EnvironmentFile=/home/ubuntu/PersonalGarminAICoach/.env
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
ExecStart=/home/ubuntu/PersonalGarminAICoach/.venv/bin/python -m streamlit run /home/ubuntu/PersonalGarminAICoach/web/app.py --server.port=8080 --server.address=0.0.0.0
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
